# ops/systemd/impress.service
[Unit]
Description=Impress CMS (API + SPA)
After=network.target
Wants=network.target

[Service]
Type=simple
User=impress
Group=impress
WorkingDirectory=/opt/impress/backend/current
EnvironmentFile=-/opt/impress/backend/.env
ExecStart=/opt/impress/backend/current/blotting-api-latest
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/impress/data /opt/impress/uploads

[Install]
WantedBy=multi-user.target
